# dist/linux/init/counterpartyd.conf
#!upstart
description "CounterParty Daemon"

# When to start the service
start on startup

# When to stop the service
stop on shutdown

# Automatically restart process if crashed
respawn

env USER=counterpartyd
env LOG_DIR=/var/log/counterpartyd
exec start-stop-daemon --start --background --make-pidfile --pidfile /var/run/counterpartyd.pid --chuid $USER --exec /usr/local/bin/counterpartyd
